# autostart/qubes-dbus-proxy.desktop
[Desktop Entry]
Icon=qubes
Name=Qubes D-Bus Proxy
Categories=System
Exec=python3 -mqubesdbus.proxy
Terminal=false
NoDisplay=true
Type=Application
